# etc/nginx/periodic/weekly/acme-client
#!/bin/sh

hosts=""

for host in $hosts; do
        acme-client -Nnmv $host && renew=1
done

[ "$renew" = 1 ] && rc-service nginx reload